# launchd/com.masa.chat-cards-hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.masa.chat-cards-hourly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/usr/bin/python3 ~/.claude/scripts/chat_card_extract.py &amp;&amp; ~/.claude/scripts/vault-prompt-runner.sh "$HOME/Documents/Obsidian Vault/00_General/prompts/scheduled/chat-cards-hourly.md" &amp;&amp; /usr/bin/python3 ~/.claude/scripts/chat_card_apply.py "$HOME/Documents/Obsidian Vault/00_Inbox/chat-reports/chat-cards-hourly-result.md"</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>18</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/com.masa.chat-cards-hourly.log</string>
  <key>StandardErrorPath</key><string>/tmp/com.masa.chat-cards-hourly.err</string>
</dict>
</plist>
