# launchd/com.masa.wiki-daily-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.masa.wiki-daily-ingest</string>
  <!--
    host: masaaki (MASA.local) — masa-2 では load しない
      （このジョブは 1 台のホスト機に固定する。別 Mac (masa-2) では launchctl load しないこと。
       result frontmatter の generated_host / apply-log jsonl で実走ホストを突合し重複実行を検知する）

    第二の脳 v3 = 日次全自動・ハブ追記型 ingest（決定 2026-07-06「無人AI書込はハブ追記型」）。
    毎日 8:47 発火（news 収集 8:07 の後）。2 段:
      ① vault-prompt-runner.sh が headless claude -p で .raw 新着を採点し
         wiki/meta/wiki-daily-ingest-result.md（固定名・上書き）へ patch(json フェンス1個)+人間サマリーを出力
      ② wiki_ingest_apply.py が result の patch を柵越しに機械適用（ハブ `## AI追記` 節へ追記のみ）
    runner が非ゼロ終了なら && で applier は実行されない（壊れた result を適用しない）。
    適用先・上限・秘密検査・冪等・flock は applier 側にハードコード（vault md から任意経路を作らせない）。

    TCC 対策 (2026-07-17 確定): bash 直実行は ~/Documents を TCC に無音拒否されるため、
    VaultJobs.app（許可ダイアログ取得済み app bundle）経由で実行する。ジョブ本体の分岐は
    EnvironmentVariables の CLAUDE_JOB → ~/.claude/scripts/vault-job-dispatch.sh。
    ⚠️ VaultJobs.app を再ビルドすると TCC 許可が無効化される（再生成は build-vaultjobs-app.sh）。

    ⚠️ 正本はステージング (~/.claude/launchd/)。activate:
        cp ~/.claude/launchd/com.masa.wiki-daily-ingest.plist ~/Library/LaunchAgents/
        launchctl unload ~/Library/LaunchAgents/com.masa.wiki-daily-ingest.plist 2>/dev/null
        launchctl load   ~/Library/LaunchAgents/com.masa.wiki-daily-ingest.plist
    手動テスト: launchctl start com.masa.wiki-daily-ingest
    停止:       launchctl unload ~/Library/LaunchAgents/com.masa.wiki-daily-ingest.plist
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/VaultJobs.app/Contents/MacOS/applet</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CLAUDE_JOB</key><string>wiki</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/masaaki</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>47</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/masaaki/.claude/state/wiki-daily-ingest.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/masaaki/.claude/state/wiki-daily-ingest.err.log</string>
</dict>
</plist>
