# launchd/com.masa.eng-vocab-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.masa.eng-vocab-weekly</string>
  <!-- TCC対策 2026-07-17: bash直実行は ~/Documents を無音拒否されるため VaultJobs.app 経由。
       本体分岐は CLAUDE_JOB → vault-job-dispatch.sh。app再ビルドは許可無効化に注意 -->
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/VaultJobs.app/Contents/MacOS/applet</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CLAUDE_JOB</key><string>eng</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>23</integer></dict>
  <key>StandardOutPath</key><string>/tmp/com.masa.eng-vocab-weekly.log</string>
  <key>StandardErrorPath</key><string>/tmp/com.masa.eng-vocab-weekly.err</string>
</dict>
</plist>
