# launchd/com.masa.chat-cards-hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.masa.chat-cards-hourly</string>
  <!-- TCC対策 2026-07-17: bash直実行だとDocuments(vault)を読めずFDA付与のGUI操作が必要になる。
       app bundle経由ならOS標準の「書類フォルダへのアクセス許可」ダイアログが自動で出て、
       ユーザーは「許可」1クリックで済む（chain実体は ChatCards.app 内の do shell script） -->
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/ChatCards.app/Contents/MacOS/applet</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>18</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>18</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/com.masa.chat-cards-hourly.log</string>
  <key>StandardErrorPath</key><string>/tmp/com.masa.chat-cards-hourly.err</string>
</dict>
</plist>
